# deployment/crontab
# Run notification check every minute
* * * * * root /app/run.sh >> /var/log/cron.log 2>&1
# Empty line required at the end
